# com.things3.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.things3.sync</string>
    
    <!-- Disable standard output/error to prevent log file creation -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
    
    <!-- Run the script through bash with proper environment -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$HOME/Development/Things3TodaySync" &amp;&amp; source "$HOME/.bash_profile" &amp;&amp; "./sync_things.sh" &gt;/dev/null 2&gt;&amp;1</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>$HOME/Development/Things3TodaySync</string>

    <!-- Set environment variables if needed -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>$HOME</string>
        <key>USER</key>
        <string>$USER</string>
        <key>SHELL</key>
        <string>/bin/bash</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep the process running, restart if it dies -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- Ensure clean process management -->
    <key>AbandonProcessGroup</key>
    <true/>
    
    <!-- Run with low priority -->
    <key>LowPriorityIO</key>
    <true/>
    
    <!-- Don't let the system kill this process -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist> 