# server/services/com.things3.google-tasks-interceptor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.things3.google-tasks-interceptor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-c</string>
        <string>cd /Users/markofir/Development/Things3-Android-Companion-App-Detailed/server && source ../venv/bin/activate && export $(cat .env | xargs) && python automated_workflow.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/markofir/Development/Things3-Android-Companion-App-Detailed/server</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/markofir/Development/Things3-Android-Companion-App-Detailed/server/logs/google-tasks-interceptor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/markofir/Development/Things3-Android-Companion-App-Detailed/server/logs/google-tasks-interceptor.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist> 